# ubuntu20-fs/var/lib/dpkg/info/fprintd.list
/.
/etc
/etc/fprintd.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/fprintd.service
/usr
/usr/bin
/usr/bin/fprintd-delete
/usr/bin/fprintd-enroll
/usr/bin/fprintd-list
/usr/bin/fprintd-verify
/usr/libexec
/usr/libexec/fprintd
/usr/share
/usr/share/dbus-1
/usr/share/dbus-1/interfaces
/usr/share/dbus-1/interfaces/net.reactivated.Fprint.Device.xml
/usr/share/dbus-1/interfaces/net.reactivated.Fprint.Manager.xml
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/net.reactivated.Fprint.service
/usr/share/dbus-1/system.d
/usr/share/dbus-1/system.d/net.reactivated.Fprint.conf
/usr/share/doc
/usr/share/doc/fprintd
/usr/share/doc/fprintd/NEWS.Debian.gz
/usr/share/doc/fprintd/changelog.Debian.gz
/usr/share/doc/fprintd/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/fprintd.1.gz
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/net.reactivated.fprint.device.policy
/var
/var/lib
/var/lib/fprint
/usr/share/man/man1/fprintd-delete.1.gz
/usr/share/man/man1/fprintd-enroll.1.gz
/usr/share/man/man1/fprintd-list.1.gz
/usr/share/man/man1/fprintd-verify.1.gz
