# ubuntu20-fs/var/lib/dpkg/info/fwupd.list
/.
/etc
/etc/fwupd
/etc/fwupd/daemon.conf
/etc/fwupd/redfish.conf
/etc/fwupd/remotes.d
/etc/fwupd/remotes.d/lvfs-testing.conf
/etc/fwupd/remotes.d/lvfs.conf
/etc/fwupd/remotes.d/vendor-directory.conf
/etc/fwupd/remotes.d/vendor.conf
/etc/fwupd/thunderbolt.conf
/etc/fwupd/uefi_capsule.conf
/etc/grub.d
/etc/grub.d/35_fwupd
/etc/pki
/etc/pki/fwupd
/etc/pki/fwupd/GPG-KEY-Linux-Foundation-Firmware
/etc/pki/fwupd/GPG-KEY-Linux-Vendor-Firmware-Service
/etc/pki/fwupd/LVFS-CA.pem
/etc/pki/fwupd-metadata
/etc/pki/fwupd-metadata/GPG-KEY-Linux-Foundation-Metadata
/etc/pki/fwupd-metadata/GPG-KEY-Linux-Vendor-Firmware-Service
/etc/pki/fwupd-metadata/LVFS-CA.pem
/etc/update-motd.d
/etc/update-motd.d/85-fwupd
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/fwupd-offline-update.service
/lib/systemd/system/fwupd-refresh.service
/lib/systemd/system/fwupd-refresh.timer
/lib/systemd/system/fwupd.service
/lib/systemd/system/system-update.target.wants
/lib/systemd/system-preset
/lib/systemd/system-shutdown
/lib/systemd/system-shutdown/fwupd.shutdown
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/90-fwupd-devices.rules
/usr
/usr/bin
/usr/bin/dbxtool
/usr/bin/dfu-tool
/usr/bin/fwupdate
/usr/bin/fwupdmgr
/usr/bin/fwupdtool
/usr/lib
/usr/lib/aarch64-linux-gnu
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_acpi_dmar.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_acpi_facp.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_acpi_phat.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_amt.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_analogix.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_ata.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_bcm57xx.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_bios.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_ccgx.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_colorhug.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_cpu.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_cros_ec.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_dell_dock.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_dfu.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_dfu_csr.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_ebitdo.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_elanfp.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_elantp.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_emmc.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_ep963x.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_fastboot.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_fresco_pd.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_genesys.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_goodixmoc.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_hailuck.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_iommu.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_jabra.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_lenovo_thinklmi.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_linux_lockdown.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_linux_sleep.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_linux_swap.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_linux_tainted.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_logind.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_logitech_hidpp.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_modem_manager.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_mtd.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_nitrokey.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_nordic_hid.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_nvme.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_optionrom.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_parade_lspcon.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_pci_bcr.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_pci_mei.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_pixart_rf.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_realtek_mst.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_redfish.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_rts54hid.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_rts54hub.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_scsi.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_steelseries.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_superio.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_synaptics_cape.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_synaptics_cxaudio.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_synaptics_mst.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_synaptics_prometheus.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_synaptics_rmi.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_system76_launch.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_thelio_io.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_thunderbolt.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_tpm.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_uefi_capsule.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_uefi_dbx.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_uefi_pk.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_uefi_recovery.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_uf2.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_upower.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_usi_dock.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_vli.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_wacom_raw.so
/usr/lib/aarch64-linux-gnu/fwupd-plugins-5/libfu_plugin_wacom_usb.so
/usr/libexec
/usr/libexec/fwupd
/usr/libexec/fwupd/fwupd
/usr/libexec/fwupd/fwupdoffline
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/fwupdagent
/usr/share/bash-completion/completions/fwupdmgr
/usr/share/bash-completion/completions/fwupdtool
/usr/share/dbus-1
/usr/share/dbus-1/interfaces
/usr/share/dbus-1/interfaces/org.freedesktop.fwupd.xml
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.fwupd.service
/usr/share/dbus-1/system.d
/usr/share/dbus-1/system.d/org.freedesktop.fwupd.conf
/usr/share/doc
/usr/share/doc/fwupd
/usr/share/doc/fwupd/copyright
/usr/share/fish
/usr/share/fish/vendor_completions.d
/usr/share/fish/vendor_completions.d/fwupdmgr.fish
/usr/share/fwupd
/usr/share/fwupd/add_capsule_header.py
/usr/share/fwupd/firmware_packager.py
/usr/share/fwupd/install_dell_bios_exe.py
/usr/share/fwupd/metainfo
/usr/share/fwupd/metainfo/org.freedesktop.fwupd.remotes.lvfs-testing.metainfo.xml
/usr/share/fwupd/metainfo/org.freedesktop.fwupd.remotes.lvfs.metainfo.xml
/usr/share/fwupd/quirks.d
/usr/share/fwupd/quirks.d/analogix.quirk
/usr/share/fwupd/quirks.d/ata.quirk
/usr/share/fwupd/quirks.d/bcm57xx.quirk
/usr/share/fwupd/quirks.d/ccgx-ids.quirk
/usr/share/fwupd/quirks.d/ccgx.quirk
/usr/share/fwupd/quirks.d/cfi.quirk
/usr/share/fwupd/quirks.d/colorhug.quirk
/usr/share/fwupd/quirks.d/cpu.quirk
/usr/share/fwupd/quirks.d/cros-ec.quirk
/usr/share/fwupd/quirks.d/dell-dock.quirk
/usr/share/fwupd/quirks.d/dfu-csr.quirk
/usr/share/fwupd/quirks.d/dfu.quirk
/usr/share/fwupd/quirks.d/ebitdo.quirk
/usr/share/fwupd/quirks.d/elanfp.quirk
/usr/share/fwupd/quirks.d/elantp.quirk
/usr/share/fwupd/quirks.d/emmc.quirk
/usr/share/fwupd/quirks.d/ep963x.quirk
/usr/share/fwupd/quirks.d/fastboot.quirk
/usr/share/fwupd/quirks.d/fresco-pd.quirk
/usr/share/fwupd/quirks.d/genesys.quirk
/usr/share/fwupd/quirks.d/goodixmoc.quirk
/usr/share/fwupd/quirks.d/hailuck.quirk
/usr/share/fwupd/quirks.d/iommu.quirk
/usr/share/fwupd/quirks.d/jabra.quirk
/usr/share/fwupd/quirks.d/logitech-hidpp.quirk
/usr/share/fwupd/quirks.d/modem-manager.quirk
/usr/share/fwupd/quirks.d/mtd.quirk
/usr/share/fwupd/quirks.d/nitrokey.quirk
/usr/share/fwupd/quirks.d/nordic-hid.quirk
/usr/share/fwupd/quirks.d/nvme.quirk
/usr/share/fwupd/quirks.d/optionrom.quirk
/usr/share/fwupd/quirks.d/parade-lspcon.quirk
/usr/share/fwupd/quirks.d/pci-bcr.quirk
/usr/share/fwupd/quirks.d/pci-mei.quirk
/usr/share/fwupd/quirks.d/pixart-rf.quirk
/usr/share/fwupd/quirks.d/power.quirk
/usr/share/fwupd/quirks.d/realtek-mst.quirk
/usr/share/fwupd/quirks.d/redfish.quirk
/usr/share/fwupd/quirks.d/rts54hid.quirk
/usr/share/fwupd/quirks.d/rts54hub.quirk
/usr/share/fwupd/quirks.d/scsi.quirk
/usr/share/fwupd/quirks.d/steelseries.quirk
/usr/share/fwupd/quirks.d/superio.quirk
/usr/share/fwupd/quirks.d/synaptics-cape.quirk
/usr/share/fwupd/quirks.d/synaptics-cxaudio.quirk
/usr/share/fwupd/quirks.d/synaptics-mst.quirk
/usr/share/fwupd/quirks.d/synaptics-prometheus.quirk
/usr/share/fwupd/quirks.d/synaptics-rmi.quirk
/usr/share/fwupd/quirks.d/system76-launch.quirk
/usr/share/fwupd/quirks.d/thelio-io.quirk
/usr/share/fwupd/quirks.d/thunderbolt.quirk
/usr/share/fwupd/quirks.d/tpm.quirk
/usr/share/fwupd/quirks.d/uefi-capsule.quirk
/usr/share/fwupd/quirks.d/uefi-dbx.quirk
/usr/share/fwupd/quirks.d/uefi-recovery.quirk
/usr/share/fwupd/quirks.d/uf2.quirk
/usr/share/fwupd/quirks.d/usi-dock.quirk
/usr/share/fwupd/quirks.d/vli-bizlink.quirk
/usr/share/fwupd/quirks.d/vli-dell.quirk
/usr/share/fwupd/quirks.d/vli-goodway.quirk
/usr/share/fwupd/quirks.d/vli-hyper.quirk
/usr/share/fwupd/quirks.d/vli-lenovo.quirk
/usr/share/fwupd/quirks.d/vli-samsung.quirk
/usr/share/fwupd/quirks.d/wacom-raw.quirk
/usr/share/fwupd/quirks.d/wacom-usb.quirk
/usr/share/fwupd/remotes.d
/usr/share/fwupd/remotes.d/vendor
/usr/share/fwupd/remotes.d/vendor/firmware
/usr/share/fwupd/remotes.d/vendor/firmware/README.md
/usr/share/fwupd/simple_client.py
/usr/share/fwupd/uefi-capsule-ux.tar.xz
/usr/share/icons
/usr/share/icons/hicolor
/usr/share/icons/hicolor/scalable
/usr/share/icons/hicolor/scalable/apps
/usr/share/icons/hicolor/scalable/apps/org.freedesktop.fwupd.svg
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/fwupd
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/dbxtool.1.gz
/usr/share/man/man1/dfu-tool.1.gz
/usr/share/man/man1/fwupdagent.1.gz
/usr/share/man/man1/fwupdate.1.gz
/usr/share/man/man1/fwupdmgr.1.gz
/usr/share/man/man1/fwupdtool.1.gz
/usr/share/metainfo
/usr/share/metainfo/org.freedesktop.fwupd.metainfo.xml
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.fwupd.policy
/usr/share/polkit-1/rules.d
/usr/share/polkit-1/rules.d/org.freedesktop.fwupd.rules
/var
/var/cache
/var/cache/app-info
/var/cache/app-info/xmls
/var/lib
/var/lib/polkit-1
/var/lib/polkit-1/localauthority
/var/lib/polkit-1/localauthority/10-vendor.d
/var/lib/polkit-1/localauthority/10-vendor.d/fwupd.pkla
/lib/systemd/system/system-update.target.wants/fwupd-offline-update.service
/usr/bin/fwupdagent
/usr/share/doc/fwupd/changelog.Debian.gz
