# ubuntu20-fs/var/lib/dpkg/info/pulseaudio.conffiles
/etc/dbus-1/system.d/pulseaudio-system.conf
/etc/init.d/pulseaudio-enable-autospawn
/etc/pulse/daemon.conf
/etc/pulse/default.pa
/etc/pulse/system.pa
/etc/xdg/autostart/pulseaudio.desktop
