# ubuntu20-fs/var/lib/dpkg/info/switcheroo-control.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/net.hadess.SwitcherooControl.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/switcheroo-control.service
/lib/udev
/lib/udev/hwdb.d
/lib/udev/hwdb.d/30-pci-intel-gpu.hwdb
/usr
/usr/libexec
/usr/libexec/switcheroo-control
/usr/share
/usr/share/doc
/usr/share/doc/switcheroo-control
/usr/share/doc/switcheroo-control/NEWS.gz
/usr/share/doc/switcheroo-control/README.md
/usr/share/doc/switcheroo-control/changelog.Debian.gz
/usr/share/doc/switcheroo-control/copyright
